# deploy/flask-admin.service
[Unit]
Description=Elysian Admin Panel (Flask)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/dubai-portal/extractor
Environment="PATH=/var/www/dubai-portal/extractor/venv/bin"
EnvironmentFile=/var/www/dubai-portal/extractor/.env
ExecStart=/var/www/dubai-portal/extractor/venv/bin/gunicorn \
    --workers 1 \
    --bind 127.0.0.1:5001 \
    --timeout 120 \
    --log-level info \
    --access-logfile /var/log/flask-admin-access.log \
    --error-logfile /var/log/flask-admin-error.log \
    admin:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
